# conf.d/rabbitmq-start.conf
[program:rabbitmq-start]
command=/usr/local/bin/rabbitmq-start
