# conf.d/rabbitmq-cluster.conf
[program:rabbitmq-cluster]
command=/usr/local/bin/rabbitmq-cluster
